# etc/profile.d/baseimage.sh
# Aliases
alias cd="HOME=~/workspace cd"
alias cp="cp -i"
alias ll="ls -l --color=auto"
alias mv="mv -i"
alias pip="pip3"
alias python="python3"
alias rm="rm -i"

# Environment variables
export EDITOR=nano

# File mode creation mask
umask 0077
